# Sources/Edge/Resources/com.edgeos.helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.edgeos.helper</string>
	<key>ProgramArguments</key>
	<array>
		<string>__EDGE_HELPER_PATH__</string>
		<string>--foreground</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>__HOME_DIR__/Library/Logs/edge-helper.log</string>
	<key>StandardErrorPath</key>
	<string>__HOME_DIR__/Library/Logs/edge-helper.log</string>
	<key>WorkingDirectory</key>
	<string>__HOME_DIR__</string>
	<key>ProcessType</key>
	<string>Background</string>
</dict>
</plist>
